# code-server.service
[Unit]
After=network.target

[Service]
User=ubuntu
Group=ubuntu
ExecStart=/usr/local/bin/codeboot.sh

[Install]
WantedBy=default.target
